# cpack-package-control-extra-scripts/postinst
echo "POSTinst $@"

systemctl start dbus-org.maui-accounts-dbus-daemon.service
systemctl enable dbus-org.maui-accounts-dbus-daemon.service
