# cpack-package-control-extra-scripts/preinst
echo "PREinst $@"

systemctl is-active --quiet dbus-org.maui-accounts-dbus-daemon.service && systemctl stop dbus-org.maui-accounts-dbus-daemon.service
systemctl is-active --quiet dbus-org.maui-accounts-dbus-daemon.service && systemctl disable dbus-org.maui-accounts-dbus-daemon.service

exit 0