# cpack-package-control-extra-scripts/prerm
echo "PRErm $@"

systemctl stop dbus-org.maui-accounts-dbus-daemon.service
systemctl disable dbus-org.maui-accounts-dbus-daemon.service
